# ansible/roles/bomberman/files/bomberman-deploy.path
[Unit]
Description=Bomberman Deploy Watcher

[Path]
PathChanged=/webapps/bomberman

[Install]
WantedBy=multi-user.target
